# veriscope_ta_dashboard/docker/8.0/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord.pid
minfds=10240


[program:php-fpm]
command=/usr/sbin/php-fpm8.4 -F -R
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true

[program:websockets]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan websockets:serve  --host=0.0.0.0 --port=6001
user=serviceuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord/websocket.log

[program:worker]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan queue:work
user=serviceuser
autostart=false
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord/worker.log

[program:horizon]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan horizon
user=serviceuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord/horizon.log

[program:scheduler]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan schedule:work
user=serviceuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord/scheduler.log

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord/cron.log
